# flight-direct/flight-completion/etc/profile.d/bash-completion.sh
#
# Sources the completion script (if it exists)
#

if [[ -f "$FL_ROOT"/scripts/bash_completion.sh ]]; then
  source "$FL_ROOT"/scripts/bash_completion.sh
fi

